# version-1/installer/crontab
#
# MARK FOR IF
## Run Speedtest.net
*/10 * * * * root bash /opt/speedpi/run_speedtest.sh
## Run Pingtest
*/1 * * * * root cd /opt/speedpi; /usr/bin/python3 /opt/speedpi/pinger.py
